# src/app/docker/base/etc/s6-overlay/s6-rc.d/dev-server/run
#!/command/with-contenv sh
exec pnpm --dir /var/www dev
